# deploy/quadlets/racklab-horizon-app.container
[Unit]
Description=RackLab Horizon (app pool: provider + notifications, no Podman socket)
Documentation=https://github.com/cyberbalsa/racklab
Requires=racklab-postgres.service racklab-redis.service racklab-plugin-bootstrap.service
After=racklab-postgres.service racklab-redis.service racklab-plugin-bootstrap.service
PartOf=racklab-runtime.target

[Container]
Image=ghcr.io/cyberbalsa/racklab/horizon:main
ContainerName=racklab-horizon-app
Network=racklab.network
EnvironmentFile=/etc/racklab/racklab.env
Environment=RACKLAB_HORIZON_POOL_GROUP=app
Volume=/etc/racklab:/etc/racklab:ro,Z
Volume=/var/lib/racklab/storage:/var/www/html/storage:Z
Volume=/var/lib/racklab/plugins:/var/lib/racklab/plugins:ro,Z
Exec=php artisan horizon
StopSignal=SIGTERM
StopTimeout=3700
NoNewPrivileges=true
Pull=newer

[Service]
Restart=always
TimeoutStartSec=300
TimeoutStopSec=3730

[Install]
WantedBy=racklab-runtime.target
