# deploy/quadlets/racklab-reverb.container
[Unit]
Description=RackLab Reverb websocket server
Documentation=https://github.com/cyberbalsa/racklab
Requires=racklab-postgres.service racklab-redis.service racklab-plugin-bootstrap.service
After=racklab-postgres.service racklab-redis.service racklab-plugin-bootstrap.service
PartOf=racklab-runtime.target

[Container]
Image=ghcr.io/cyberbalsa/racklab/reverb:main
ContainerName=racklab-reverb
Network=racklab.network
EnvironmentFile=/etc/racklab/racklab.env
Volume=/etc/racklab:/etc/racklab:ro,Z
Volume=/var/lib/racklab/storage:/var/www/html/storage:Z
Volume=/var/lib/racklab/plugins:/var/lib/racklab/plugins:Z
PublishPort=0.0.0.0:8080:8080
Exec=php artisan reverb:start --host=0.0.0.0 --port=8080
HealthCmd=php artisan about --only=environment
HealthInterval=30s
HealthRetries=3
HealthTimeout=10s
NoNewPrivileges=true
Pull=newer

[Service]
Restart=always
TimeoutStartSec=300
TimeoutStopSec=120

[Install]
WantedBy=racklab-runtime.target
