# deploy/quadlets/racklab-scheduler-reconciler@.container
[Unit]
Description=RackLab scheduler and reconciler %i
Documentation=https://github.com/cyberbalsa/racklab
Requires=racklab-postgres.service racklab-redis.service racklab-plugin-bootstrap.service
After=racklab-postgres.service racklab-redis.service racklab-plugin-bootstrap.service
PartOf=racklab-runtime.target

[Container]
Image=ghcr.io/cyberbalsa/racklab/scheduler-reconciler:main
ContainerName=racklab-scheduler-reconciler-%i
Network=racklab.network
EnvironmentFile=/etc/racklab/racklab.env
Environment=RACKLAB_WORKER_POOL=scheduler-reconciler
Environment=RACKLAB_WORKER_REPLICA=%i
Volume=/etc/racklab:/etc/racklab:ro,Z
Volume=/var/lib/racklab/storage:/var/www/html/storage:Z
Volume=/var/lib/racklab/plugins:/var/lib/racklab/plugins:ro,Z
# Laravel scheduler: dispatches the reconcile/expire/drift/reap jobs onto Redis
# at their scheduled cadence (see bootstrap/app.php withSchedule). Horizon
# workers (maintenance → app pool, cleanup → runner pool) process them. The
# scheduler only enqueues, so it no longer needs the Podman socket.
Exec=php artisan schedule:work
StopSignal=SIGTERM
StopTimeout=60
NoNewPrivileges=true
Pull=newer

[Service]
Restart=always
TimeoutStartSec=300
TimeoutStopSec=90

[Install]
WantedBy=racklab-runtime.target
DefaultInstance=1
